# ops/ai.openclaw.court-watch-nightly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.court-watch-nightly</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/varysoc/.openclaw/workspace/court-watch/.venv/bin/python</string>
      <string>/Users/varysoc/.openclaw/workspace/court-watch/scripts/nightly_capture_and_report.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/varysoc/.openclaw/workspace/court-watch</string>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>22</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/varysoc/.openclaw/workspace/court-watch/logs/nightly.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/varysoc/.openclaw/workspace/court-watch/logs/nightly.stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
  </dict>
</plist>
